# bitbake-recipe/files/geophone-rotate-heatmap.cron
# Rotate the heatmap images every 15 minutes
*/15 * * * *	root	PATH="/usr/bin:$PATH" /usr/bin/geophone-rotate-heatmap
